# DEBIAN/postrm
#!/bin/bash
update-rc.d simplydo remove

